# cinema/templates/systemd/cinema.service
[Unit]
Description=Mise à jour des programmes de cinéma

[Service]
Type=oneshot
ExecStart=/opt/cinema.py
